# debian/max-control.postrm
#!/bin/sh


if [ "$1" = "purge" ]; then

   if [ -x "/etc/init.d/apache2" ]; then
    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d apache2 reload || true
    else
        /etc/init.d/apache2 reload || true
    fi
  fi

fi

#DEBHELPER#
